# systemd/delta-algo.service
[Unit]
Description=Delta Algo Trading Framework (Streamlit Dashboard)
Documentation=https://github.com/your-org/delta-algo-framework
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rdpuser

# ── Paths ────────────────────────────────────────────────────────────────────
WorkingDirectory=/home/rdpuser/delta-algo-framework

# EnvironmentFile loads variables into the process environment.
# NOTE: systemd does NOT expand ${VAR} tokens inside ExecStart=.
# Variable substitution is handled inside start-streamlit.sh instead.
EnvironmentFile=/home/rdpuser/delta-algo-framework/.env

# ── Startup ──────────────────────────────────────────────────────────────────
# start-streamlit.sh sources .env itself, resolves DASHBOARD_PORT,
# then exec's the venv Streamlit binary with the concrete port number.
ExecStart=/bin/bash /home/rdpuser/delta-algo-framework/deploy/start-streamlit.sh

# ── Process supervision ──────────────────────────────────────────────────────
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=20

# ── Logging (view with: journalctl -u delta-algo -f) ────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=delta-algo

[Install]
WantedBy=multi-user.target
